# kanakku-admin-server.service
[Unit]
Description=Kanakku Admin Server for System Administration
After=network.target
Wants=kanakku.service

[Service]
Type=simple
User=kanakku
Group=kanakku
WorkingDirectory=/opt/kanakku/adminserver
Environment="PATH=/opt/kanakku/adminserver/venv/bin"
Environment="PYTHONPATH=/opt/kanakku/adminserver"
EnvironmentFile=/opt/kanakku/.env
ExecStart=/opt/kanakku/adminserver/venv/bin/python admin_server.py
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal

# Security enhancements
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/kanakku/logs /var/log/kanakku /opt/kanakku/adminserver/logs
NoNewPrivileges=true
PrivateDevices=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

# Network access for SSH connections
PrivateNetwork=false

[Install]
WantedBy=multi-user.target
