# kanakku-admin-server.service
[Unit]
Description=Kanakku Admin Server for System Administration
After=network.target
Wants=kanakku.service

[Service]
Type=simple
User=kanakku
Group=kanakku
WorkingDirectory=/opt/kanakku/adminserver
Environment="PATH=/opt/kanakku/adminserver/venv/bin"
Environment="PYTHONPATH=/opt/kanakku/adminserver"
Environment="KANAKKU_DEPLOY_HOST=localhost"
Environment="KANAKKU_DEPLOY_USER=kanakku"
Environment="KANAKKU_SSH_KEY_PATH=/home/kanakku/.ssh/id_rsa"
EnvironmentFile=/opt/kanakku/.env
ExecStart=/opt/kanakku/adminserver/venv/bin/python admin_server.py
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal

# Security enhancements (minimal restrictions for admin functionality)
PrivateTmp=false
ProtectSystem=false
NoNewPrivileges=false
PrivateDevices=false
ProtectHome=false
PrivateNetwork=false
PrivateMounts=false

# Explicitly allow access to required paths
ReadWritePaths=/opt/kanakku
ReadWritePaths=/home/kanakku
ReadWritePaths=/var/log
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target
